# template/memdeflate@.service
[Service]
Type=oneshot
ExecStart=/opt/memdeflate/memdeflate --vmname=%i
